# pkg/debian/archivebox.service
# The archivebox user/group and /var/lib/archivebox directory are created by
# postinstall.sh (which runs after dpkg unpacks the package contents).

[Unit]
Description=ArchiveBox Web Archiving Server
After=network.target

[Service]
Type=simple
User=archivebox
Group=archivebox
WorkingDirectory=/var/lib/archivebox
Environment="PATH=/opt/archivebox/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStartPre=/usr/bin/archivebox init
ExecStart=/usr/bin/archivebox server 0.0.0.0:8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
